# supervisord.conf
[supervisord]
nodaemon=true

[program:grpc_gateway]
command=/usr/bin/argusagent
autostart=true
autorestart=true
stdout_logfile=/var/log/argusagent.log
stderr_logfile=/var/log/argusagent_err.log
